# Library/LaunchAgents/de.patwoz.crontab.backuphome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>de.patwoz.crontab.backuphome</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/patwoz/.bin/backup-home.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>20</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/de.patwoz.crontab.backuphome.err</string>

  <key>StandardOutPath</key>
  <string>/tmp/de.patwoz.crontab.backuphome.out</string>
</dict>
</plist>
